# container-files/etc/supervisord.d/npc.conf
[program:npc]
user=root
command=/config/npc/npc 
process_name=npc
stdout_logfile=/var/log/npc.stdout.log
stderr_logfile=/var/log/npc.stderr.log
autorestart=true
autostart=true
